# config/supervisord/deployment.dev.conf
[supervisord]

[program:app]
command=pipenv run python manage.py runserver 0.0.0.0:8080
environment=HOME="/home/app",USER="app",PIPENV_DONT_LOAD_ENV="true"
directory=/opt/service/app
user=app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0